# systemd/agecalc-rss.service
[Unit]
Description=AgeCalc RSS ingest and blog generation
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=oneshot
User=agecalc
Group=www-data
WorkingDirectory=/srv/apps/agecalc
Environment="PATH=/srv/apps/agecalc/.micromamba/envs/agecalc/bin"
EnvironmentFile=/srv/apps/agecalc/.env.rss
# Required vars live in /srv/apps/agecalc/.env.rss. See .env.rss.example.
ExecStart=/srv/apps/agecalc/.micromamba/envs/agecalc/bin/python /srv/apps/agecalc/scripts/rss_blog_scheduler.py run --limit 1 --status draft --provider openai --model gpt-4.1-mini
